# infra/systemd/openclaw-runtime-from-1password.service
[Unit]
Description=OpenClaw Gateway (runtime secrets from 1Password)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/var/lib/openclaw

# Render ephemeral runtime env from 1Password into tmpfs (/run).
ExecStartPre=/usr/local/bin/openclaw-env-from-1password /run/openclaw/env
EnvironmentFile=/run/openclaw/env

# Keep all secret references in config as ${VAR} placeholders.
ExecStart=/usr/local/bin/openclaw gateway run --port 18789

Restart=always
RestartSec=5
TimeoutStopSec=30

# Baseline hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/openclaw /var/lib/openclaw

[Install]
WantedBy=multi-user.target
