# infra/systemd/openclaw-clawdoc.service
[Unit]
Description=OpenClaw claw-doc Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/.openclaw-clawdoc

# Render runtime secrets from 1Password into tmpfs before start.
ExecStartPre=/usr/local/bin/render-env-from-1password.sh /run/openclaw-clawdoc/env

# Runtime env file (ephemeral, not checked into git)
EnvironmentFile=/run/openclaw-clawdoc/env

# Start in foreground (systemd supervises process lifecycle)
ExecStart=/home/%i/.npm-global/bin/openclaw gateway run --port 18890

Restart=always
RestartSec=5
TimeoutStopSec=30

# Hardening baseline
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/openclaw-clawdoc /home/%i/.openclaw-clawdoc

[Install]
WantedBy=multi-user.target
